# pialert/rootfs/etc/s6-overlay/s6-rc.d/install-pialert/up
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start install-pialert service
# ==============================================================================
#bashio::log.info "Starting install-pialert..."

/etc/s6-overlay/scripts/install-pialert.sh