# sl2file/rootfs/etc/s6-overlay/s6-rc.d/syslog-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start syslog-server service
# ==============================================================================
bashio::log.green "Starting up syslog-server..."
if $(bashio::config 'debug'); then
    bashio::log.level "debug"
fi

# Setting PYTHONUNBUFFERED to a non empty value ensures that the python output
# is sent straight to the container log.
# This also ensures that no partial output is held in a buffer somewhere and
# never written in case the python application crashes.

export PYTHONUNBUFFERED=1

bashio::log.green "Saving configuration options to file..."

CONFIG_PATH=/config.py
{
    echo "DEBUG_MODE = \"$(bashio::config 'debug')\"";
    echo "CUSTOM_HOST = \"$(bashio::config 'custom_host')\"";
    echo "CUSTOM_PORT = $(bashio::addon.port '514/udp')";
    echo "CUSTOM_LOG_FILE = \"$(bashio::config 'custom_log_file')\"";
} > "${CONFIG_PATH}"

bashio::log.green "Running syslog-server..."

python3 /sl2file.py

bashio::log.info "syslog-server stopped"
exit 0
