# pialert/rootfs/etc/s6-overlay/s6-rc.d/lighttpd/run
#!/command/with-contenv bash

if [ "${PH_VERBOSE:-0}" -gt 0 ] ; then
    set -x ;
fi

lighttpd -D -f /etc/lighttpd/lighttpd.conf
