# pi.alert/rootfs/etc/s6-overlay/s6-rc.d/pialert-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start pialert-server service
# ==============================================================================
bashio::log.green "Starting up pialert-server..."
if $(bashio::config 'debug'); then
    bashio::log.level "debug"
fi

# Setting PYTHONUNBUFFERED to a non empty value ensures that the python output
# is sent straight to the container log.
# This also ensures that no partial output is held in a buffer somewhere and
# never written in case the python application crashes.

export PYTHONUNBUFFERED=1

bashio::log.green "Saving configuration options to file..."

CONFIG_PATH=/config.py
{
    echo "DEBUG_MODE = \"$(bashio::config 'debug')\"";
} > "${CONFIG_PATH}"

bashio::log.green "Running pialert-server..."

#python3 /pi.alert.py

bashio::log.info "pialert-server stopped"
exit 0
